# debian/dispmanx_vnc.socket
[Unit]
Description=dispmanx_vnc Socket for Per-Connection Servers

[Socket]
ListenStream=5901
Accept=true

[Install]
WantedBy=sockets.target
